# systemd/telegram-control-bot.service
[Unit]
Description=Telegram Control Bot
After=network.target telegram-backend.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-platform/bots/control_bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target