# systemd/telegram-streaming-bot.service
[Unit]
Description=Telegram Streaming Bot
After=network.target telegram-backend.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-platform/bots/streaming_bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target